# llms-web.service
[Unit]
Description=LLMS Web Server
After=network.target

[Service]
Type=simple
User=house
WorkingDirectory=/home/house/llms
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
